# deploy/ansible/roles/uwsgi/templates/uwsgi.conf
description "daemon for uwsgi server"

start on startup
stop on shutdown
respawn

script
  uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
end script
